# startEmus.conf
start on runlevel [2345]
stop on shutdown

chdir /home/ubuntu/snesio.SMB/hqnclient/

script
	./runCpus.sh 2> /home/ubuntu/snesio.SMB/error.log
end script
