# startEmus.conf
start on runlevel [2345]
stop on shutdown

chdir /usr/local/bin/snesio.SMB/hqnclient/

script
	./runCpus.sh 2> /usr/local/bin/snesio.SMB/error.log
end script
